# ops/systemd/munich-weekly-backup.service
[Unit]
Description=Munich Weekly production backup
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target
OnFailure=munich-weekly-alert@%n.service

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/sbin/munich-weekly-backup.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
