# ops/systemd/munich-weekly-alert@.service
[Unit]
Description=Send Munich Weekly operations alert for %i

[Service]
Type=oneshot
User=root
Group=root
PassEnvironment=ALERT_ENV
ExecStart=/usr/local/sbin/munich-weekly-notify-ops.sh %i
